# roles/admin_bastion/templates/bastion.service.j2
[Unit]
Description=Update the bastion host configuration

[Service]
Type=simple
ExecStart=/usr/bin/su - {{ admin_username }} -c bastion.sh

[Install]
WantedBy=multi-user.target
